# data/check-update.timer.in
[Unit]
Description=Periodic check for new OS version

[Timer]
OnBootSec=5m
OnUnitActiveSec=1h
AccuracySec=20m

[Install]
WantedBy=multi-user.target
